# support/vagrant/service/notebook/run
#!/bin/sh

# move stderr to stdout
exec 2>&1

exec su vagrant -c 'ambry ui notebook -w -H 0.0.0.0 -p8081'
